# configs/releng/airootfs/etc/udev/rules.d/90-pcspkr-beep.rules
# Add write access to the PC speaker for the "beep" group
ACTION=="add", SUBSYSTEM=="input", ATTRS{name}=="PC Speaker", ENV{DEVNAME}!="", RUN+="/usr/bin/setfacl -m u:beep:w '$env{DEVNAME}'"
